# deploy/systemd/osm-rangefind-nightly.service
[Unit]
Description=OSM Rangefind nightly indexing window
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/osm-rangefind-index
Environment=INDEX_LOG_FILE=logs/nightly.log
ExecStart=/srv/osm-rangefind-index/scripts/nightly.sh --deadline 06:15

